# debian/logstash.conf
[program:logstash]
command=/usr/share/logstash/logstash.sh
user=logstash
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=QUIT