# debian/logstash.conf
[program:logstash]
command=python /usr/share/logstash/logstash.py
user=logstash
process_name=%(program_name)s
autostart=true
autorestart=true
startsecs=3
stopsignal=KILL